# docker_template/backend/laravel/autostart.conf
[supervisord]
nodaemon=true

[program:laravel-migrate]
command=php /var/www/html/artisan migrate
directory=/var/www/html
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/migrate.log

[program:laravel-webserver]
command=php /var/www/html/artisan serve --host=0.0.0.0 --port=8000
directory=/var/www/html
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/webserver.log

[program:laravel-queue]
command=php /var/www/html/artisan queue:work
directory=/var/www/html
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queues.log